# weather/systemd/mk-piclock-api.service.d/weather-source.conf
[Service]
ReadWritePaths=/var/lib/mk-piclock-weather
